# systemd/conscio.service
[Unit]
Description=Conscio autonomous conscious agent service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/conscio
Environment=CONSCIO_CONFIG=/home/conscio/.conscio/config.toml
ExecStart=/opt/conscio/.venv/bin/conscio service start
Restart=always
RestartSec=5
# Hardening — chosen to not break the agent's bash/execute_code tools.
# Excluded on purpose: ProtectSystem/ProtectHome (the agent legitimately writes
# its working directory and ~/.conscio), MemoryDenyWriteExecute (numpy/cffi).
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectClock=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
LimitNOFILE=65536
# Exit code 3 = corrupt state database (see docs/runbooks/db-locked-or-corrupted.md);
# restarting cannot help and would crash-loop.
RestartPreventExitStatus=3
User=conscio
Group=conscio

[Install]
WantedBy=multi-user.target
